# bin/config/profile.d/merge.cmd
"%CMDER_ROOT%/php/php.exe" "%CMDER_ROOT%/script/merge"